# service/rd-monitor-daemon.service
[Unit]
Description=RD Monitor - daemon to auto-fix waiting torrents
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Projets_Github/rd-monitor
ExecStart=%h/Projets_Github/rd-monitor/scripts/run_service.sh --daemon --persist %h/Projets_Github/rd-monitor/data/auto_fix_state.db --results %h/Projets_Github/rd-monitor/data/auto_fix_results.jsonl --pause 1.5 --max-per-cycle 200 --cycle-interval 3600
Restart=on-failure
RestartSec=30
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
